# config/launchd/com.canon.meta_continuous_scan.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.n6-architecture.meta_continuous_scan</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/ghost/core/n6-architecture/bin/n6_meta</string>
    <string>continuous-scan</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/ghost/core/n6-architecture</string>
  <key>StartInterval</key>
  <integer>43200</integer>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/tmp/n6_meta_scan.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/n6_meta_scan.err.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>N6</key>
    <string>/Users/ghost/core/n6-architecture</string>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
